# deploy/service/espnet.service.in
[Unit]
Description=ESPnet serving service
After=network.target

[Service]
ExecStart=make run-service
WorkingDirectory=${PWD}
StandardOutput=journal
StandardError=journal
Restart=always
User=${USER}

[Install]
WantedBy=multi-user.target